# files/etc/supervisor.d/lexicon.ini
[program:lexicon]
command = bash /app/lexicon.sh
startretries= 1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
